# config/auto_publish_hexo.conf
[program:auto_hexo]
directory = /root/auto_publish
command = node deploy.js
autostart = yes
autorestart = yes

