# config/auto_publish_hexo.conf
[program:hexo]
directory = /root
command = node deploy.js
autostart = yes
autorestart = yes

