# config/hexo.conf
[program:hexo]
directory = /root/blog
command = /usr/bin/hexo server -i 0.0.0.0
autostart = yes
autorestart = yes

