# config/sshd.conf
[program:openssh]
command=/usr/sbin/sshd
